# btrfs/btrfs-scrub.service
[Unit]
Description=btrfs scrub

[Service]
Type=simple
ExecStart=/usr/local/bin/btrfs-scrub
