# btrfs/btrfs-scrub.timer
[Unit]
Description=Runs btrfs scrub on all discs monthly

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target